# aws/systemd/careflowai-backend.service
[Unit]
Description=CareFlowAI FastAPI Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/careflowai/backend
Environment="PATH=/opt/careflowai/backend/venv/bin"
ExecStart=/opt/careflowai/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
